# vm-agent/echo-vue-agent.service
[Unit]
Description=SocialPlug VM Agent — by Illy Robotic Instruments
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/socialplug-agent
ExecStart=/home/ubuntu/socialplug-agent/venv/bin/python agent.py --vm-id VM_ID_HERE --bucket livepay-petition.firebasestorage.app
Restart=always
RestartSec=10
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
